# deploy/csr.service
[Unit]
Description=CSR IRIT Prototype API
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/csr-irit/backend
Environment="PATH=/opt/csr-irit/venv/bin"
Environment="SECRET_KEY=CHANGE_ME_IN_PRODUCTION"
Environment="TOKEN_TTL_MIN=360"
ExecStart=/opt/csr-irit/venv/bin/gunicorn \
    --workers 3 \
    --bind 127.0.0.1:8000 \
    --timeout 120 \
    --access-logfile /var/log/csr-irit/access.log \
    --error-logfile /var/log/csr-irit/error.log \
    app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
